# ops/systemd/vibenews-api.service
# VibeNews private API (loopback-only; exposed to the tailnet solely via Tailscale Serve).
# Authored artifact — NOT installed or started by the implementation pass. Install/enable is a
# reviewed operator step. EnvironmentFile holds the mode-600 secrets that are never committed.
[Unit]
Description=VibeNews private API (loopback)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vibenews
Group=vibenews
WorkingDirectory=/opt/vibenews
EnvironmentFile=/etc/vibenews/.env.server.local
ExecStart=/usr/bin/node --import tsx server/src/bin/api.ts
Restart=on-failure
RestartSec=3
# Hardening: loopback bind is enforced in-process; these narrow the blast radius further.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vibenews-dev
StateDirectory=vibenews-dev
CapabilityBoundingSet=
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
